# docker/supervisor-service.ini
[supervisord]
nodaemon = true
pidfile = /var/run/supervisord.pid
logfile = /var/log/supervisor.log

[program:haproxy]
autostart = true
autorestart = true
command = haproxy -f /etc/haproxy/haproxy.cfg

# The following config enables haproxy reload through supervisorctl

[unix_http_server]
file = /run/supervisor.sock   

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///run/supervisor.sock
